# roles/ovs-dpdk-switch/files/openvswitch-nonetwork.service.d/10-depend_on_dpdk.conf
[Unit]
Requires=dpdk.service
After=dpdk.service
