# roles/docker/files/docker.service.d/service-overrides.conf
[Service]
EnvironmentFile=-/etc/default/docker
EnvironmentFile=-/etc/default/docker-storage
ExecStart=
ExecStart=/usr/bin/dockerd $OPTIONS \
        $DOCKER_STORAGE_OPTIONS
