# files/supervisor/faucet.redcables.wand.nz/bird6.conf
[program:bird6]
command=/usr/sbin/bird6 -f
