# files/supervisor/faucet.redcables.wand.nz/bird.conf
[program:bird]
command=/usr/sbin/bird -f
