# files/supervisor/faucet.redcables.wand.nz/faucet.conf
[program:faucet]
command=ryu-manager --config-file=/etc/ryu/ryu.conf faucet.faucet
